# ansible/roles/bambi-openvpn-team-server/files/setup-team-openvpn.service
[Unit]
Description=Setup OpenVPN team server

[Service]
ExecStart=/opt/setup-team-openvpn.sh

[Install]
WantedBy=multi-user.target
