# ansible/roles/enomoloch-get-pcaps/files/get_tcp_dump_messages.service
[Unit]
Description=getting enodump

[Service]
Restart=always
TimeoutStartSec=300

ExecStart=/services/sh get_messages.sh

[Install]
WantedBy=multi-user.targe